# com.sushant.locus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sushant.locus</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Locus.app/Contents/MacOS/locus</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardErrorPath</key>
    <string>/tmp/locus.err</string>
    
    <key>StandardOutPath</key>
    <string>/tmp/locus.out</string>
    
    <key>WorkingDirectory</key>
    <string>/Applications/Locus.app/Contents/MacOS</string>
    
    <!-- Optional: Restart if the process crashes -->
    <key>CrashTimeout</key>
    <integer>10</integer>
    
    <!-- Optional: Throttle restarts to avoid rapid cycling -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>